# crontab
30 02 * * *  /Users/jerome/bin/backup_to_external_hd.sh >> /Users/jerome/conf/log/cron.log 2>&1
